# supervisord.conf
[supervisord]
nodaemon=true

[program:line_bot]
command=/myapp/line_botx.linux
stdout_logfile=/myapp/line_botx.log
stdout_logfile_maxbytes=100


